# systemd/.config/systemd/mpd.socket
[Socket]
ListenStream=7700

[Install]
WantedBy=default.target
