# systemd/.config/systemd/system-maint.target
[Unit]
Description=Runs System Maintenance daily

[Timer]
OnUnitActiveSec=1day

[Install]
WantedBy=default.target
